# home/Library/LaunchAgents/com.prateek.tartelet-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Launch Tartelet at login on a runner host and keep it running. Paired with the
  startVirtualMachinesOnLaunch=true key in the managed dk.shape.Tartelet.plist, so
  login -> running ephemeral runners is unattended. Gated in .chezmoiignore on the
  tartelet cask, so it only lands on hosts that install Tartelet (homelab).

  Loads at the next login automatically; to load immediately after apply:
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.prateek.tartelet-runner.plist

  Runner-guest LAN isolation is not set here. Tartelet cannot pass networking flags
  to `tart run`, so the tart wrapper at /opt/homebrew/bin/tart adds --net-softnet
  instead (installed by run_after_18-tartelet-tart-softnet-wrapper). See the runbook.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.prateek.tartelet-runner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Tartelet.app/Contents/MacOS/Tartelet</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
</dict>
</plist>
